# files/rk35xx/etc/uci-defaults/99_lan_dhcp_ignore
uci get dhcp.lan 1>/dev/null 2>/dev/null && {
uci batch <<EOF
set dhcp.lan.ignore='1'
commit dhcp
EOF
}